# radioProject.conf
start on startup
exec forever start /full/path/to/app.js
